# logstash-offline/roles/logstash-ansible/templates/logstash.service.j2
[Unit]
Description=Logstash Service

[Service]
User={{ ls_user }}
ExecStart=/etc/logstash/bin/logstash -f /etc/logstash/conf.d

[Install]
Alias=logstash.service
